# projects/ota_basic_example/ota_package_build/zoho_iot_ota/zoho_iot_ota.service
[Unit]
Description= zoho_iot_ota

[Service]
ExecStart= /tmp/zoho_iot_ota/zupdate.sh
WorkingDirectory= /tmp/zoho_iot_ota
User=root
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target