# projects/ota_deb_example/deb_package_build/zoho_iot_ota/etc/systemd/system/zoho_iot_deb.service
[Unit]
Description= zoho_iot_deb

[Service]
ExecStart= /usr/sbin/ota_deb
Restart=always
User=root
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
