# projects/ota_deb_example/deb_package_build/zoho_iot_ota/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable zoho_iot_deb.service
systemctl start zoho_iot_deb.service

exit 0
